# setting/supervisor.conf
[supervisord]
logfile=/home/runner/work/knight_server/knight_server/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/home/runner/work/knight_server/knight_server/logs/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=runner
childlogdir=/home/runner/work/knight_server/knight_server/logs

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:gate_8000]
command=/usr/bin/python -m services.gate.gate_server --port 8000
process_name=gate_8000
directory=/home/runner/work/knight_server
autostart=false
autorestart=true
startsecs=5
startretries=3
user=runner
redirect_stderr=false
stdout_logfile=/home/runner/work/knight_server/knight_server/logs/gate_8000.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/home/runner/work/knight_server/knight_server/logs/gate_8000.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/home/runner/work/knight_server",SERVICE_NAME="gate",SERVICE_PORT="8000",SERVICE_INSTANCE="gate-8000",ENVIRONMENT="development",SERVICE_WORKERS="1",SERVICE_ACCESS_LOG="True",SERVICE_AUTO_RELOAD="False",SERVICE_WEBSOCKET_MAX_SIZE="1048576",SERVICE_WEBSOCKET_PING_INTERVAL="10",SERVICE_WEBSOCKET_PING_TIMEOUT="30",SERVICE_REQUEST_TIMEOUT="60",SERVICE_RESPONSE_TIMEOUT="60",SERVICE_KEEP_ALIVE_TIMEOUT="5",SERVICE_RATE_LIMIT="{'enabled': True, 'max_requests': 100, 'window': 60}"
priority=999
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:logic_70052]
command=/usr/bin/python -m services.logic.logic_server --port 70052
process_name=logic_70052
directory=/home/runner/work/knight_server
autostart=false
autorestart=true
startsecs=5
startretries=3
user=runner
redirect_stderr=false
stdout_logfile=/home/runner/work/knight_server/knight_server/logs/logic_70052.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/home/runner/work/knight_server/knight_server/logs/logic_70052.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/home/runner/work/knight_server",SERVICE_NAME="logic",SERVICE_PORT="70052",SERVICE_INSTANCE="logic-70052",ENVIRONMENT="development",SERVICE_WORKERS="4",SERVICE_MAX_CONCURRENT_REQUESTS="1000",SERVICE_MAX_RECEIVE_MESSAGE_LENGTH="4194304",SERVICE_MAX_SEND_MESSAGE_LENGTH="4194304",SERVICE_REQUEST_TIMEOUT="30",SERVICE_CACHE_TTL="300"
priority=989
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:chat_70051]
command=/usr/bin/python -m services.chat.chat_server --port 70051
process_name=chat_70051
directory=/home/runner/work/knight_server
autostart=false
autorestart=true
startsecs=5
startretries=3
user=runner
redirect_stderr=false
stdout_logfile=/home/runner/work/knight_server/knight_server/logs/chat_70051.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/home/runner/work/knight_server/knight_server/logs/chat_70051.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/home/runner/work/knight_server",SERVICE_NAME="chat",SERVICE_PORT="70051",SERVICE_INSTANCE="chat-70051",ENVIRONMENT="development",SERVICE_WORKERS="2",SERVICE_MAX_MESSAGE_LENGTH="500",SERVICE_MAX_ROOM_USERS="100",SERVICE_MESSAGE_RATE_LIMIT="10",SERVICE_HISTORY_MESSAGE_COUNT="100",SERVICE_HISTORY_RETENTION_DAYS="7"
priority=979
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[program:fight_70053]
command=/usr/bin/python -m services.fight.fight_server --port 70053
process_name=fight_70053
directory=/home/runner/work/knight_server
autostart=false
autorestart=true
startsecs=5
startretries=3
user=runner
redirect_stderr=false
stdout_logfile=/home/runner/work/knight_server/knight_server/logs/fight_70053.out.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/home/runner/work/knight_server/knight_server/logs/fight_70053.err.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
environment=PYTHONPATH="/home/runner/work/knight_server",SERVICE_NAME="fight",SERVICE_PORT="70053",SERVICE_INSTANCE="fight-70053",ENVIRONMENT="development",SERVICE_WORKERS="2",SERVICE_MAX_BATTLE_TIME="300",SERVICE_BATTLE_TICK_RATE="0.1",SERVICE_MAX_PLAYERS_PER_BATTLE="10",SERVICE_BATTLE_POOL_SIZE="100",SERVICE_ENABLE_BATTLE_RECORD="True"
priority=969
stopasgroup=true
killasgroup=true
stopsignal=TERM
stopwaitsecs=10

[group:all]
programs=gate_8000,logic_70052,chat_70051,fight_70053
priority=999

[group:gate]
programs=gate_8000
priority=999

[group:logic]
programs=logic_70052
priority=999

[group:chat]
programs=chat_70051
priority=999

[group:fight]
programs=fight_70053
priority=999